# templates/etc/systemd/system/docker-prune.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Prune unused Docker images
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune --force --volumes
